# root/etc/services.d/activemq/run
#!/usr/bin/with-contenv bash

cd ${ACTIVEMQ_HOME} || exit 1;

s6-setuidgid alpine exec activemq console
